# deploy/vector-owner-control.timer.example
[Unit]
Description=Poll owner commands for Vector

[Timer]
OnBootSec=15s
OnUnitActiveSec=15s
AccuracySec=1s
Unit=vector-owner-control.service

[Install]
WantedBy=timers.target
