# deploy/vector-owner-control.timer.example
[Unit]
Description=Poll owner commands for Vector

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true
Unit=vector-owner-control.service

[Install]
WantedBy=timers.target
