# deploy/vector-ops-approval.timer.example
[Unit]
Description=Poll WhatsApp approvals for Vector Ops

[Timer]
OnBootSec=20s
OnUnitActiveSec=20s
AccuracySec=1s
Unit=vector-ops-approval.service

[Install]
WantedBy=timers.target
